# tools/publish/systemd/publish@.timer
[Unit]
Description=Rebuild the public site for instance %i once daily

[Timer]
OnCalendar=*-*-* 06:00:00
Persistent=true
RandomizedDelaySec=10min

[Install]
WantedBy=timers.target
