# tools/audit-tap/systemd/audit-tap@.timer
[Unit]
Description=Run Continuant audit tap for instance %i every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
